# overlays/18-firmware-updates/usr/lib/systemd/system/update-initramfs.service
[Unit]
Description=InitramFS Update Service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/neon/update_initramfs.sh
