# overlays/46-audio-receiver/etc/udev/rules.d/99-bluetooth-audio.rules
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="*Bluetooth Audio*", RUN+="/usr/local/bin/a2dp.sh %k"