# overlays/29-ovos-shell/usr/lib/systemd/system/gui-shell.service
[Unit]
Description=Neon GUI

[Service]
User=1000
WorkingDirectory=/opt/neon
EnvironmentFile=/etc/neon/gui-shell.conf
ExecStartPre=-/usr/bin/kbuildsycoca5
ExecStart=ovos-shell
Restart=on-failure

[Install]
WantedBy=multi-user.target
